# deploy/tilinx-web.service
[Unit]
Description=TilinX Web Dashboard
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/tilinx/website
EnvironmentFile=/opt/tilinx/.env
ExecStart=/usr/bin/python3 app.py
Restart=always
RestartSec=5
StandardOutput=append:/var/log/tilinx/web.log
StandardError=append:/var/log/tilinx/web.log

[Install]
WantedBy=multi-user.target
